# crontab
* * * * * /bin/bash /home/ross/git/cm4/add_fish.sh
